# cronjob/crontab
SHELL=/bin/bash
#
01 00 * * * root cd /covid/lab && /covid/job >/proc/1/fd/1 2>/proc/1/fd/2
#
